# mock_hostname/ld.so.preload
/usr/lib64/libdl.so.2
/usr/local/lib64/libmock_hostname.so
